# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/local/openresty/nginx/sbin/nginx
