# openresty/openresty.conf
[program:openresty]
command=/opt/openresty/nginx/sbin/nginx
autorestart=true
